# configurations/app/app.ini
[program:oneline]
directory=/app
command=bash -c 'python app.py'
stdout_logfile=/logs/oneline_supervisor.log
stderr_logfile=/logs/oneline_supervisor_error.log
autostart=true
autorestart=true
startsecs=10
killasgroup=true
startretries=10
stopasgroup=true